# pkgbuild/g13d.install
post_install() {
    systemctl enable --now g13.target
    systemctl enable --now g13-resume.service
}

pre_remove() {
    systemctl disable --now g13.target
    systemctl disable --now g13-resume.service
}