# systemd/g13-clock.service
[Unit]
Description=Logitech G13 clock
After=g13.service
Requires=g13.service

[Service]
Type=simple
ExecStart=/usr/local/bin/g13-clock -o /run/g13d/g13-0

[Install]
WantedBy=g13.service

