# systemd/g13.target
[Unit]
Description=Logitech G13 Hotplug Target
Requires=g13.service
After=usb.target

[Install]
WantedBy=multi-user.target