# host/systemd/system/g13d.target
[Unit]
Description=Logitech G13 Hotplug Target
Requires=g13d.service
After=usb.target

[Install]
WantedBy=multi-user.target